# roles/common/templates/node_exporter.conf.j2
description "exporter for system stats"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec node_exporter -web.listen-address={{ node_exporter_listen_address }}
